# docker/supervisord-sandbox.conf
[program:nginx]
command=nginx -g "daemon off;"
redirect_stderr=true

[program:gunicorn]
command=/usr/local/bin/gunicorn --config /app/docker/gunicorn.py ietf.wsgi
directory=/app
redirect_stderr=true